# systemd/user/barinfo.service
[Unit]
Description=Swaybar info

[Service]
ExecStart=%h/.local/bin/barinfo
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
